# build/crontab
0 3 * * 0 /home/homeassistant/update-ha.sh